# config/cron/crontab
SHELL={{command -v fish}}

@reboot clean_temp_uvars
