# config/cron/crontab
SHELL={{command -v fish}}

@reboot clean_toggles
